# configs/.config/profile.d/01-environment-variables.sh
export QT_QPA_PLATFORMTHEME="qt5ct"

export CL_RESET=$(echo -en '\033[0m')
export CL_RED=$(echo -en '\033[00;31m')
export CL_GREEN=$(echo -en '\033[00;32m')
export CL_YELLOW=$(echo -en '\033[00;33m')
export CL_BLUE=$(echo -en '\033[00;34m')
export CL_MAGENTA=$(echo -en '\033[00;35m')
export CL_PURPLE=$(echo -en '\033[00;35m')
export CL_CYAN=$(echo -en '\033[00;36m')
export CL_LGRAY=$(echo -en '\033[00;37m')
export CL_LRED=$(echo -en '\033[01;31m')
export CL_LGREEN=$(echo -en '\033[01;32m')
export CL_LYELLOW=$(echo -en '\033[01;33m')
export CL_LBLUE=$(echo -en '\033[01;34m')
export CL_LMAGENTA=$(echo -en '\033[01;35m')
export CL_LPURPLE=$(echo -en '\033[01;35m')
export CL_LCYAN=$(echo -en '\033[01;36m')
export CL_WHITE=$(echo -en '\033[01;37m')